# deploy/server/sharepaste-server.service
[Unit]
Description=SharePaste gRPC Server
After=network-online.target postgresql.service redis.service
Wants=network-online.target

[Service]
Type=simple
User=sharepaste
Group=sharepaste
WorkingDirectory=/opt/sharepaste
EnvironmentFile=/etc/sharepaste/server.env
ExecStart=/usr/bin/env bash /opt/sharepaste/scripts/start-server-prod.sh
Restart=always
RestartSec=3
TimeoutStopSec=20
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
